# debian/vyatta-cfg-vpn.postinst.in
#!/bin/bash

set -e

prefix=@prefix@
exec_prefix=@exec_prefix@
sysconfdir=@sysconfdir@
bindir=@bindir@
sbindir=@sbindir@

systemctl enable vyos-xfrmwatchd.service

exit 0
